# chef/letsencrypt/crontab

# Automatic renewal LetsEncrypt cert every 1.month, at: '2am'
0 2 1 * * certbot renew --post-hook "service nginx reload" --webroot-path /var/www/letsencrypt >> /var/log/letsencrypt_renew.log